# devel/etc/catkin/profile.d/50.movo_network_config.bash
. "/home/andrewandjasmine/movo_ws/src/movo_network/movo_network_config.bash"
