# devel/etc/catkin/profile.d/50.movo_config.bash
. "/home/andrewandjasmine/movo_ws/src/movo_common/movo_config/movo_config.bash"
